# backup_20250312_221412/backend-backup-20250311-034535/conf/harmonic-universe.service
[Unit]
Description=Harmonic Universe Gunicorn Daemon
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/backend
Environment="PATH=/path/to/your/venv/bin"
Environment="PYTHONPATH=/path/to/your/backend"
ExecStart=/path/to/your/venv/bin/gunicorn --config /path/to/your/backend/conf/gunicorn.conf.py wsgi:application

[Install]
WantedBy=multi-user.target
